# board/rpi0/rootfs-overlay/etc/sv/klogd/run
#!/bin/runit-sv

sv start syslogd

SV_ARGS="-n"
