# board/rpi0/rootfs-overlay/etc/sv/hakomari-displayd/run
#!/bin/runit-sv

DEVICE="/dev/${SV_NAME_PARAM}"
SV_ARGS="/dev/i2c"
RUN_DIR="/var/run/hakomari"
GPIO_DIR="/sys/class/gpio"
DEVICES_DIR="/sys/devices"
SANDBOX_MOUNTS="\
	${GPIO_DIR}:${GPIO_DIR}:rw \
	${RUN_DIR}:${RUN_DIR}:rw \
	${DEVICE}:${SV_ARGS}:rw \
	${DEVICES_DIR}:${DEVICES_DIR}:rw \
"
SANDBOX_USER=hakomari-displayd
SANDBOX_GROUP=hakomari-daemon
SV_EXEC="exec sandbox-exec"

mkdir -p /var/run/hakomari
chown :hakomari-daemon /var/run/hakomari
chmod u=rwx,g=rwx,o= /var/run/hakomari
