# board/rpi0/rootfs-overlay/etc/sv/hakomari-dispatcherd/run
#!/bin/runit-sv

DEVICE="/dev/${SV_NAME_PARAM}"
SV_ARGS="/dev/serial"
RUN_DIR="/var/run/hakomari"
SANDBOX_MOUNTS="\
	${DEVICE}:${SV_ARGS}:rw \
	${RUN_DIR}:${RUN_DIR}:rw \
"
SANDBOX_USER=hakomari-dispatcherd
SANDBOX_GROUP=hakomari-daemon
SV_EXEC="exec sandbox-exec"

chown ${SANDBOX_USER}:${SANDBOX_GROUP} ${DEVICE}
