# board/rpi0/rootfs-overlay/etc/sv/dropbear/run
#!/bin/runit-sv

mkdir -p /var/run/dropbear

SV_ARGS="-R -F -m -j -k -B"

[ "${DEV_MODE}" != "1" ] && sv down ${SV_NAME} || true
