# board/rpi0/rootfs-overlay/etc/sv/watchdog/run
#!/bin/runit-sv

DEVICE="/dev/watchdog"
REBOOT_TIMEOUT=5
RESET_INTERVAL=3

[ ! -e "${DEVICE}" ] && sv down ${SV_NAME} || true

SV_ARGS="-F -T ${REBOOT_TIMEOUT} -t ${RESET_INTERVAL} ${DEVICE}"
