# board/rpi0/rootfs-overlay/etc/sv/watchdog/run
#!/bin/runit-sandboxed-sv

DEVICE="/dev/watchdog"
REBOOT_TIMEOUT=5
RESET_INTERVAL=3
SV_ARGS="-F -T ${REBOOT_TIMEOUT} -t ${RESET_INTERVAL} ${DEVICE}"
SANDBOX_ARGS="--mount ${DEVICE}:/dev/watchdog --read-only"

[ ! -e "${DEVICE}" ] && sv down ${SV_NAME} || true
