# board/rpi0/rootfs-overlay/etc/sv/udhcpd/run
#!/bin/runit-sv

touch /var/run/udhcpd.leases

SV_ARGS="-f /etc/udhcpd.conf"
