# board/rpi0/rootfs-overlay/etc/sv/udhcpd/run
#!/bin/runit-sv

SV_ARGS="-f /etc/udhcpd.conf"
SV_EXEC="exec sandbox-exec"
SANDBOX_MOUNTS="\
	/etc/udhcpd.conf:/etc/udhcpd.conf \
	/var/run/udhcpd.leases:/run/udhcpd.leases:rw \
"
SANDBOX_ARGS="--network"

touch /var/run/udhcpd.leases

[ "${DEV_MODE}" != "1" ] && sv down ${SV_NAME} || true
