# board/rpi0/rootfs-overlay/etc/sv/syslogd/run
#!/bin/runit-sv

SV_ARGS="-n -O /var/log/messages"
