# board/rpi0/rootfs-overlay/etc/sv/getty/run
#!/bin/runit-sv

[ ! -e "/dev/${SV_NAME_PARAM}" ] && sv down ${SV_NAME} || true

SV_ARGS="-L ${SV_NAME_PARAM} 115200 vt100"
